# root/etc/s6-overlay/s6-rc.d/py-comicsviewer/run
#!/usr/bin/with-contenv bash

cd /py-comicsviewer
exec python3 web.py